# wongle/software/etc/init/radish-wongle-sitter.conf
start on stopped rc RUNLEVEL=[2345]
stop on runlevel [016]

respawn
exec /usr/local/radish/bin/sitter